# services/redmine/data/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/etc/init.d/mysql start -D

[program:redmine]
command=/bin/bash -c "cd /var/www/redmine-3.1.0 && bundle exec rails server webrick -e production -b 0.0.0.0"
