# services/zimbra/data/supervisord.conf
[supervisord]
nodaemon=true

[program:scriptstart]
command=/bin/bash -c "su zimbra && zmcontrol start"
