# services/openfire/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/etc/init.d/mysql start -D

[program:openfire]
command=/etc/init.d/openfire start -D
